# infra/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[supervisord]
user=root
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
childlogdir=/var/log/supervisor
pidfile=/var/run/supervisord.pid
loglevel=info
environment=ANONYMIZED_TELEMETRY="false",CHROME_DEBUGGING_PORT="9222",CHROME_DEBUGGING_HOST="localhost",CHROME_PERSISTENT_SESSION="true",DISPLAY=:99

; ============================================================
; DISPLAY
; ============================================================
[program:xvfb]
command=Xvfb :99 -screen 0 1600x900x24 -ac +extension GLX +render -noreset
autorestart=true
stdout_logfile=/dev/null
stderr_logfile=/dev/null
priority=100
startsecs=3
stopsignal=TERM
stopwaitsecs=10

; ============================================================
; VNC — no password, no nginx
; ============================================================
[program:x11vnc]
command=bash -c "sleep 3 && exec x11vnc -display :99 -forever -shared -nopw -rfbport 5901 -o /var/log/x11vnc.log"
autorestart=true
stdout_logfile=/dev/null
stderr_logfile=/dev/null
priority=200
startretries=10
startsecs=5
stopsignal=TERM
stopwaitsecs=10
environment=DISPLAY=":99"

[program:novnc]
command=bash -c "sleep 5 && exec python3 -m websockify --web /opt/novnc 0.0.0.0:6080 localhost:5901"
environment=PYTHONPATH="/opt/novnc/utils/websockify"
autorestart=true
stdout_logfile=/var/log/supervisor/novnc.out.log
stderr_logfile=/var/log/supervisor/novnc.err.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB
priority=300
startretries=5
startsecs=8

; ============================================================
; BROWSER
; ============================================================
[program:ninja_browser]
command=/bin/bash -c "source ~/.bashrc 2>/dev/null; cd /workspace/ninja && exec python -m browser.browser_server start --foreground"
directory=/workspace/ninja
user=root
autostart=true
autorestart=true
startsecs=5
startretries=10
environment=DISPLAY=":99"
stderr_logfile=/var/log/supervisor/ninja_browser.err.log
stdout_logfile=/var/log/supervisor/ninja_browser.out.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB

; ============================================================
; PROXY
; ============================================================
[program:psiphon]
command=/workspace/ninja/browser/psiphon-tunnel-core -config /workspace/ninja/browser/psiphon.config.json
directory=/workspace/ninja/browser
user=root
autostart=true
autorestart=true
startsecs=5
startretries=10
stderr_logfile=/var/log/supervisor/psiphon.err.log
stdout_logfile=/var/log/supervisor/psiphon.out.log
stdout_logfile_maxbytes=1MB
stderr_logfile_maxbytes=1MB

; ============================================================
; AGENT DASHBOARD
; ============================================================
[program:agent-dashboard]
command=/bin/bash -c "source ~/.bashrc 2>/dev/null; exec python3 app.py"
directory=/workspace/ninja/dashboard
user=root
autostart=true
autorestart=true
startsecs=3
startretries=5
environment=DASHBOARD_PORT="9000"
stderr_logfile=/var/log/supervisor/agent-dashboard.err.log
stdout_logfile=/var/log/supervisor/agent-dashboard.out.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB

[program:integrations-dashboard]
command=/bin/bash -c "source ~/.bashrc 2>/dev/null; exec python3 integrations_app.py"
directory=/workspace/ninja/dashboard
user=root
autostart=true
autorestart=true
startsecs=3
startretries=5
environment=INTEGRATIONS_PORT="9020",PYTHONPATH="/workspace/ninja/src/ninja:/workspace/ninja"
stderr_logfile=/var/log/supervisor/integrations-dashboard.err.log
stdout_logfile=/var/log/supervisor/integrations-dashboard.out.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB

; ============================================================
; PLATFORM SERVICES
; ============================================================
[program:http_server]
command=python /app/server.py
directory=/app
autorestart=true
stdout_logfile=/var/log/app_server.log
stderr_logfile=/var/log/app_server.log
priority=400
startretries=5
startsecs=5
stopsignal=TERM
stopwaitsecs=10

[program:browser_api]
command=python /app/browser_api.py
directory=/app
autorestart=true
stdout_logfile=/var/log/app_browser.log
stderr_logfile=/var/log/app_browser.log
priority=400
startretries=5
startsecs=5
stopsignal=TERM
stopwaitsecs=10

[program:ttyd]
command=/usr/local/bin/ttyd -i lo -W -p 2222 /bin/bash
directory=/app
autorestart=true
stdout_logfile=/dev/null
stderr_logfile=/dev/null
priority=450
startretries=5
startsecs=5
stopsignal=TERM
stopwaitsecs=10

[program:code_server]
command=/usr/bin/code-server --bind-addr 127.0.0.1:4000 --auth none --user-data-dir /opt/code-server --extensions-dir /opt/code-server/extensions
directory=/app
autorestart=true
stdout_logfile=/dev/null
stderr_logfile=/dev/null
priority=500
startretries=5
startsecs=5
stopsignal=TERM
stopwaitsecs=10

[program:ninja-mcp-servers]
directory=/opt/ninja-mcp-servers
command=/bin/bash -lc "./start.sh"
environment=PATH=/usr/local/go/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin,HOME=/root,GOPATH=/root/go,GOMODCACHE=/root/go/pkg/mod,GOCACHE=/root/.cache/go-build
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
startretries=5
stopsignal=TERM
stdout_logfile=/var/log/ninja-mcp-servers.log
stderr_logfile=/var/log/ninja-mcp-servers.err.log

[program:slack-mcp]
command=/bin/bash -c "source ~/.bashrc 2>/dev/null; exec .venv/bin/python server.py"
directory=/opt/ninja-mcp-servers/slack-mcp
user=root
autostart=true
autorestart=true
startsecs=3
startretries=3
stderr_logfile=/var/log/supervisor/slack-mcp.err.log
stdout_logfile=/var/log/supervisor/slack-mcp.out.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB
